# etc/systemd/system/ai_manager.service
[Unit]
Description=AI Manager Flask App Service
After=network.target

[Service]
User=SERVICE_USER
Group=SERVICE_GROUP
WorkingDirectory=INSTALL_DIR/ai_manager
ExecStart=INSTALL_DIR/ai_manager/venv/bin/python INSTALL_DIR/ai_manager/app.py
Restart=always
StandardOutput=append:/var/log/ai_manager.log
StandardError=append:/var/log/ai_manager.log

[Install]
WantedBy=multi-user.target
